# .config/systemd/user/udiskie.service
[Unit]
Description=handle automounting

[Service]
ExecStart=/usr/bin/udiskie -anT2
Restart=always

[Install]
WantedBy=default.target

